# system/news-delivery-startup.service
[Unit]
Description=News Delivery System Startup Service
Documentation=https://github.com/news-delivery-system
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
RemainAfterExit=yes
User=kensan
Group=kensan
WorkingDirectory=/mnt/Linux-ExHDD/news-delivery-system
Environment=PATH=/mnt/Linux-ExHDD/news-delivery-system/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/mnt/Linux-ExHDD/news-delivery-system/system/startup-news-delivery.sh
ExecReload=/mnt/Linux-ExHDD/news-delivery-system/system/startup-news-delivery.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=news-delivery-startup

# Restart policy
Restart=on-failure
RestartSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/Linux-ExHDD/news-delivery-system
ReadOnlyPaths=/mnt/Linux-ExHDD/news-delivery-system/venv

[Install]
WantedBy=multi-user.target